# docs/examples/channel_restart.service
[Unit]
Description=restart channel_run.service

[Service]
Type=oneshot
ExecStart=systemctl restart channel_run.service
