# docker/data/conf.d/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:apache2]
command=apache2.sh
